# etc/init/conflict_detector.conf
description "Git Conflict Detector"
author "cielo24"

start on startup
stop on shutdown
respawn

script
 sudo -u root php -f /var/www/detect_conflicts_process.php
end script